# templates/cpupower.service.j2
[Unit]
Description=CPU Performance Governor

[Service]
ExecStart=cpupower frequency-set -g {{ cpu_governor }}
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
